# build/etc/supervisor.d/conf.d/phoenix-proxy.conf
[program:phoenix-proxy]
directory=/phoenix/bin
command=/phoenix/bin/phoenix-proxy -c /phoenix/etc/conf/proxy.ini -rpc true
numprocs=1
autostart=true
startsecs=1
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /phoenix/var/log/supervisor/phoenix-proxy.log
