# build/etc/supervisor.d/conf.d/phoenix-web.conf
[program:phoenix-web]
directory=/phoenix/bin
command=/phoenix/bin/phoenix-web -c /phoenix/etc/conf/web.toml
numprocs=1
autostart=true
startsecs=1
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /phoenix/var/log/supervisor/phoenix-web.log
